# meta-openguard/recipes-core/systemd/files/mguard3/85-hwclock.rules
# Enforce hwclock hctosys via udev rules
KERNEL=="rtc0", RUN+="/sbin/hwclock --hctosys"
